# relay/systemd/caddy.service.d/wol-relay.conf
# Drop-in for the Debian-packaged caddy.service.
# Deploy at /etc/systemd/system/caddy.service.d/wol-relay.conf.
#
# The packaged caddy.service does not load any EnvironmentFile by
# default, so the {$VAR} placeholders in the Caddyfile would expand to
# empty strings without this override. Pointing to /etc/caddy/wol-relay.env
# keeps Caddy config templating-friendly: the public Caddyfile carries
# only variable references, while the deployment-specific values
# (LE_EMAIL, RELAY_DOMAIN, CORS_ORIGIN) live one-shot on the VM and
# are never committed to the repo.

[Service]
EnvironmentFile=/etc/caddy/wol-relay.env
